# system-specific/nixos/config/autostart/Neovim wrapper.desktop
[Desktop Entry]
Exec=nvim --listen ~/.cache/nvim/a.pipe %F
Icon=nvim
Name=Neovim wrapper
Path=
Terminal=True
Type=Application
